# etc/systemd/system/pbapi.socket
[Unit]
Description=pbapi socket

[Socket]
ListenStream=/run/pbapi.sock

[Install]
WantedBy=sockets.target
